# udev/rules.d/99-input.rules.template
SUBSYSTEM=="input", GROUP="input", MODE="0660" 
KERNEL=="input[0-9]*", RUN+="/home/${::HOMEPI_USER::}/homepi/scripts/setup-bt-audio"
