# udev/rules.d/99-input.rules
SUBSYSTEM=="input", GROUP="input", MODE="0660" 
KERNEL=="input[0-9]*", RUN+="/home/pi/homepi/scripts/setup-bt-audio"
